# build/vagrant/fundraising_app.service
[Unit]
Description=PHP Server
After=home.mount network.target

[Service]
User=vagrant
Group=vagrant
ExecStart=/usr/bin/php -S 0.0.0.0:8080 -t /vagrant/web
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target